# init/nagastd-node.conf
# MySQL Service

description     "NagAstd Node"
author          "Mark Schouten <mark@tuxis.nl>"

start on (net-device-up
          and local-filesystems
	  and runlevel [2345])
stop on runlevel [016]

respawn

exec /usr/sbin/nagastd node --raddress ::1
